# wpa/99-wpa-wl.rules
SUBSYSTEM=="net", ACTION=="move", KERNEL=="wl*", RUN+="/usr/bin/sh /root/wpa-networkd.sh"
